# backend-api/deploy/supervisor/laravel-workers.conf
; Production Supervisor configuration for Laravel queue workers.
; Copy to /etc/supervisor/conf.d/ and run: supervisorctl reread && supervisorctl update

[program:laravel-worker-rides]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pick-you/backend-api/artisan queue:work redis --queue=rides --sleep=1 --tries=3 --max-time=3600 --timeout=30
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-rides.log
stopwaitsecs=3600

[program:laravel-worker-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pick-you/backend-api/artisan queue:work redis --queue=notifications --sleep=3 --tries=3 --max-time=3600 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-notifications.log
stopwaitsecs=3600

[program:laravel-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pick-you/backend-api/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-default.log
stopwaitsecs=3600

[group:laravel-workers]
programs=laravel-worker-rides,laravel-worker-notifications,laravel-worker-default
priority=999
